# systemd/user/delta-daemon.service
[Unit]
Description=DELTA Local AI Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/delta
Environment=PYTHONUNBUFFERED=1
Environment=DELTA_DATA_DIR=%h/.local/share/delta
Environment=DELTA_DAEMON_SOCKET=%t/delta-daemon.sock
Environment=DELTA_IPC_SOCKET=%t/delta-ipc.sock
ExecStart=/usr/bin/env delta-daemon
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/delta %t

[Install]
WantedBy=default.target
