# systemd/user/delta-daemon.service
[Unit]
Description=DELTA Local AI Daemon
After=network.target

[Service]
Type=exec
WorkingDirectory=%h
Environment=PYTHONUNBUFFERED=1
Environment=DELTA_DATA_DIR=%h/.local/share/deltai
Environment=DELTA_CONFIG_DIR=%h/.config/deltai
Environment=DELTA_CACHE_DIR=%h/.cache/deltai
Environment=DELTA_SQLITE_PATH=%h/.local/share/deltai/delta.db
Environment=DELTA_DAEMON_SOCKET=%t/deltai/daemon.sock
Environment=DELTA_IPC_SOCKET=%t/deltai/ipc.sock
ExecStart=/usr/bin/env delta-daemon
Restart=on-failure
RestartSec=2
RuntimeDirectory=deltai
RuntimeDirectoryMode=0700
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/deltai %h/.config/deltai %h/.cache/deltai %t/deltai

[Install]
WantedBy=default.target
