# systemd/user/delta-daemon.service
[Unit]
Description=DELTA Local AI Daemon
After=network.target

[Service]
Type=exec
WorkingDirectory=%h
Environment=PYTHONUNBUFFERED=1
Environment=DELTA_DATA_DIR=%h/.local/state/delta
Environment=DELTA_SQLITE_PATH=%h/.local/state/delta/delta.db
Environment=DELTA_DAEMON_SOCKET=%t/delta/daemon.sock
Environment=DELTA_IPC_SOCKET=%t/delta/ipc.sock
ExecStart=/usr/bin/env delta-daemon
Restart=on-failure
RestartSec=2
RuntimeDirectory=delta
RuntimeDirectoryMode=0700
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/delta %t/delta

[Install]
WantedBy=default.target
